# cookbook/templates/default/upstart.conf.erb
#
# Cookbook Name:: guardian
# Template:: upstart.conf.erb
#
author  "Rapid7 LLC <coreservices@rapid7.com>"
description "<%= @service.name %>"

setuid <%= @service.user %>
setgid <%= @service.group %>

start on (virtual-filesystems and net-device-up IFACE!=lo)
stop on shutdown

respawn
respawn limit 5 5

console log
exec <%= @service.exec_string %>
